# deb/gosecure_0.1.0_amd64/DEBIAN/postrm
#!/bin/bash

#rm -f /usr/local/bin/gosecure /etc/init.d/gosecure /etc/rc2.d/S99gosecure

update-rc.d gosecure remove

rm -f /usr/local/bin/gosecure /etc/init.d/gosecure 

